# seafile/root/etc/supervisor/conf.d/supervisor-seafile.conf
[program:seafile]
priority=20
directory=/tmp
command=/usr/bin/seafile-start
user=root
autostart=true
autorestart=false
stdout_logfile=syslog
stderr_logfile=syslog

