# git/root/etc/supervisor/conf.d/supervisor-gitlab.conf
[program:gitlab]
priority=20
directory=/tmp
command=/setup/gitlab-start
user=root
autostart=true
autorestart=false
stdout_logfile=syslog
stderr_logfile=syslog

