# git/root/etc/supervisor/conf.d/supervisor-ssh.conf
[program:ssh]
priority=10
directory=/tmp
command=/etc/init.d/ssh start
startsecs=5
user=root
autostart=true
autorestart=false
stdout_logfile=syslog
stderr_logfile=syslog

